# meta-rpi/recipes-qt5navigation/navigation/files/navigation.service
[Unit]
Description=Qt5 Navigation app startup service

[Service]
ExecStart=/usr/bin/navigation

[Install]
WantedBy=multi-user.target
